# tdd_django/deploy_tools/gunicorn-supervisor.template.conf
[program:tdd_django]
user=USERNAME
command = /home/USERNAME/sites/SITENAME/run.sh
stdout_logfile = /home/USERNAME/sites/SITENAME/logs/supervisor.logs
redirect_stderr = true