# deploy_tools/gunicorn-supervisor.template.conf
[program:tdd_django]
user=USERNAME
command = /home/USERNAME/SITENAME/run.sh
stdout_logfile = /home/USER/SITENAME/logs/supervisor.logs
redirect_stderr = true